# backend/deploy/supervisor.conf
[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/tukangdekat/backend/artisan queue:work --sleep=3 --tries=3 --timeout=60 --queue=default,payouts,notifications
autostart=true
autorestart=true
autorestart_delay=10
startsecs=10
startretries=3
user=www-data
numprocs=3
redirect_stderr=true
stdout_logfile=/var/log/laravel-queue.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopasgroup=true
stopwaitsecs=600
priority=999
process_priority=-5
environment=PATH="/var/www/tukangdekat/backend/vendor/bin:/usr/local/bin:/usr/bin:/bin",HOME="/home/www-data",LOGNAME="www-data"
